# media-sound/music/files/musicd.initd
#!/sbin/openrc-run

depend() {
	need dbus
	use logger
}

description="My personal music daemon"

supervisor=supervise-daemon
command="/usr/bin/musicd"
