# app-admin/pass-secret-service/files/pass-secret-service.initd
#!/sbin/openrc-run

description="org.freedesktop.secrets agent for pass"

depend() {
	need dbus
}

DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
supervisor=supervise-daemon
command="/usr/bin/pass-secret-service"
